# imapfilter.cron
#!/bin/sh

/usr/bin/systemctl kill -s SIGUSR1 imapfilter.service
